# sysutils/pfSense-upgrade/files/kontrol-upgrade-stage.rc
#!/bin/sh
# PROVIDE: kontrol_upgrade_stage
# REQUIRE: LOGIN FILESYSTEMS
# BEFORE:  DAEMON

. /etc/rc.subr

name=kontrol_upgrade_stage
rcvar=kontrol_upgrade_stage_enable
command=/usr/local/sbin/Kontrol-upgrade
pidfile=/var/run/Kontrol-upgrade-stage.pid

: ${kontrol_upgrade_stage_enable:=YES}
: ${kontrol_upgrade_stage_timeout:=7200}
: ${kontrol_upgrade_stage_debug_log:=/cf/conf/upgrade_log.debug.log}

start_cmd="${name}_start"

kontrol_upgrade_stage_debug()
{
	[ "${DEBUG_ENABLED:-1}" = "1" ] || return 0
	mkdir -p "$(dirname "${kontrol_upgrade_stage_debug_log}")"
	printf '%s [DEBUG][rc.kontrol_upgrade_stage] pid=%s %s\n' "$(date -u +%FT%TZ)" "$$" "$*" >> "${kontrol_upgrade_stage_debug_log}"
}

kontrol_upgrade_stage_start()
{
	state_file="/var/db/Kontrol/upgrade/state"
	kontrol_upgrade_stage_debug "start called state_file=${state_file}"
	[ -f "${state_file}" ] || {
		kontrol_upgrade_stage_debug "state file not found; nothing to do"
		return 0
	}

	pending_stage=$(awk -F= '$1=="stage" {print $2}' "${state_file}" | tail -n1)
	kontrol_upgrade_stage_debug "pending_stage=${pending_stage:-none}"
	case "${pending_stage}" in
		stage2|stage3|stage4)
			echo "Running pending Kontrol upgrade stage: ${pending_stage}"
			kontrol_upgrade_stage_debug "executing timeout=${kontrol_upgrade_stage_timeout} command=${command} -b ${pending_stage}"
			/usr/bin/timeout "${kontrol_upgrade_stage_timeout}" \
				"${command}" -b "${pending_stage}"
			rc=$?
			kontrol_upgrade_stage_debug "stage execution rc=${rc}"
			;;
		*)
			kontrol_upgrade_stage_debug "no eligible pending stage; skipping"
			return 0
			;;
	esac
}

load_rc_config $name
run_rc_command "$1"
